# files/etc/uci-defaults/98-dhcp.config
#!/bin/sh

# 获取当前物理网口数量 (ethX / ensX 等)
iface_count=$(ip link show | grep -E "^[0-9]+: (eth|ens)[0-9]" | wc -l)

if [ "$iface_count" -eq 1 ]; then
    # 确认 dhcp.lan 配置存在
    if uci get dhcp.lan &>/dev/null; then
        echo "Detected single interface, disabling DHCP server..."
        uci set dhcp.lan.ignore='1'
        uci delete dhcp.lan.ra_slaac
        uci commit dhcp
    fi
fi

exit 0
