# files/etc/uci-defaults/99-system.config
#!/bin/sh

# 设置系统时区
uci set system.@system[0].zonename='Asia/Shanghai'
uci set system.@system[0].timezone='CST-8'
uci commit system

# 重启系统服务，使配置生效
/etc/init.d/system restart

exit 0
